# test/integration/supervisord.conf
[include]
files=/usr/local/etc/supervisord.base.conf

[program:papermill]
command=papermill /home/apowers/integration-test.ipynb /home/apowers/integration-output.ipynb
priority=100
startsecs=0
autostart=true
autorestart=false

[eventlistener:processes]
command=/home/apowers/supervisord-shutdown.sh
events=PROCESS_STATE_STOPPED, PROCESS_STATE_EXITED, PROCESS_STATE_FATAL